# etc/udev/rules.d/00hdparm.rules
ACTION=="add", SUBSYSTEM=="block", KERNEL=="sda", RUN+="/sbin/hdparm -B 255 -S 0 /dev/sda"
