# etc/systemd/system/nbfc_service.timer
[Unit]
Description=Timer for the nbfc_service service

[Timer]
OnBootSec=1min

[Install]
WantedBy=timers.target
